# Motion-Config/clear-motion-videos.timer
[Unit]
Description=Clear Motion Videos Older than One Day Hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
